# systemd/temp-hwr.timer
[Unit]
Description=temp auslesen alle 10min

[Timer]
OnActiveSec=30s
OnUnitActiveSec=5m

[Install]
WantedBy=basic.target
